# sample-package/systemd/sample-package.service
[Unit]
Description=sample-package

[Service]
Type=simple
ExecStart=/usr/bin/sample-package --daemon

[Install]
WantedBy=multi-user.target
